# supervisord.conf
[supervisord]
nodaemon=false

[program:phdaemon]
command=/bin/bash -c "/data/phdaemon"
autostart=true
autorestart=true
startsecs=5

[program:status]
command=/bin/bash -c "cat /tmp/phtunnel.pid"